# patch/imm-24.10/rc.local
uci set network.cfg030f15.macaddr="8C:DA:$(date +"%d:%H:%M:%S")"
uci commit network

exit 0
